# frp/DEBIAN/prerm
#!/bin/bash
echo "Shutting down Frp server daemon service."
ystemctl stop frps
systemctl disable frps.service

